# docker/crontab
MAILTO=eeagrants@eaudeweb.ro
PYTHONPATH=/usr/local/bin
SHELL=/bin/bash

0 4 * * *  . /root/.bashrc; /bin/import.sh >> /var/local/logs/import.log 2>&1

